# .config/systemd/user/offlineimap-oneshot@thiefhu.service.d/override.conf
[Service]
ExecStartPre=%h/tools/dbLocked.sh